# weekly_display.service
[Unit]
Description=Weekly Display Tracking Web Application
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/weekly-display
Environment=PATH=/var/www/weekly-display/venv/bin
EnvironmentFile=/var/www/weekly-display/.env
Environment=FLASK_ENV=production
ExecStart=/var/www/weekly-display/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 4 --timeout 300 --access-logfile /var/www/weekly-display/logs/access.log --error-logfile /var/www/weekly-display/logs/error.log wsgi:app
Restart=always
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weekly-display

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/www/weekly-display/uploads /var/www/weekly-display/reports /var/www/weekly-display/charts /var/www/weekly-display/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target